# debian/cron.weekly
#!/bin/sh
#
# Last modification: Wed, 13 Jun 2001 14:41:48 +0100
# man-db cron weekly

# regenerate man database
if [ -x /usr/bin/mandb ]; then
    # --pidfile /dev/null so it always starts; mandb isn't really a daemon,
    # but we want to start it like one.
    start-stop-daemon --start --pidfile /dev/null \
		      --startas /usr/lib/man-db/mandb --oknodo --chuid man \
		      >/dev/null 2>/dev/null
fi
