# debian/cron.weekly
#!/bin/sh
#
# Last modification: Fri, 18 May 2001 13:54:15 +0100
# man-db cron weekly

# regenerate man database
if [ -x /usr/bin/mandb ]; then
    nice su man -c 'mandb 2>/dev/null >/dev/null'
fi
